# deploy/com.example.knowledgehub.daily-scan.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.knowledgehub.daily-scan</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/your-user/KnowledgeHub/scripts/knowledgehub-daily-scan.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/your-user/KnowledgeHub</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>12</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/your-user</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>KNOWLEDGEHUB_CONFIG</key>
        <string>/Users/your-user/KnowledgeHub/config/knowledge-bases.yaml</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/your-user/KnowledgeHub/logs/daily-scan.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/your-user/KnowledgeHub/logs/daily-scan.err.log</string>
</dict>
</plist>
